# pkg/system-config/system-config-gnome.install
post_install() {
	systemctl enable --now NetworkManager.service
	systemctl enable --now gdm.service
}
