# pkg/system-config/system-config-xfce.install
post_install() {
	post_upgrade
}

post_upgrade() {
	# lightdm
	systemctl enable --now lightdm.service

	# networking
	systemctl enable --now systemd-networkd.service
	systemctl enable --now systemd-resolved.service
	ln --symbolic --force ../run/systemd/resolve/stub-resolv.conf /etc/resolv.conf
}
