# pkg/system-config/system-config-xfce.install
post_install() {
	systemctl enable --now lightdm.service
}
