# pkg/system-config/system-config-networkd.install
post_install() {
	systemctl enable --now systemd-networkd.service
	systemctl enable --now systemd-resolved.service
	ln --symbolic --force ../run/systemd/resolve/stub-resolv.conf /etc/resolv.conf
}
